# debian/dbinfo.cron.daily
#!/bin/sh
cd /srv/dbinfo
./mkstat --all > /var/log/dbinfo/mkstat.log 2>&1
